# units/systemd-netlogd.service.in
# SPDX-License-Identifier: LGPL-2.1+
#  This file is part of systemd.

[Unit]
Description=Journal Syslog Unicast and Multicast Daemon
Documentation=man:systemd-netlogd.conf(5)
After=network.target

[Service]
ExecStart=/lib/systemd/systemd-netlogd
PrivateTmp=yes
PrivateDevices=yes
WatchdogSec=20min

[Install]
WantedBy=multi-user.target
