# units/systemd-netlogd.service.in
# SPDX-License-Identifier: LGPL-2.1-or-later
# This file is part of systemd.

[Unit]
Description=Journal Syslog Unicast and Multicast Daemon
Documentation=man:systemd-netlogd(8)
After=network.target

[Service]
ExecStart=@PKGPREFIX@/systemd-netlogd
WatchdogSec=20min

# hardening
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
StateDirectory=systemd/journal-netlogd
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
